# systemd/check-app.service
[Unit]
Description=Проверка здоровья FastAPI
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/chack_ap.sh