# systemd/fastapi-counter.service
[Unit]
Description=FastAPI Counter App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fastapi-counter
ExecStart=/opt/fastapi-counter/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target